# src/fs/xfs/xfsprogs-dev/scrub/xfs_scrub_all.service.in
[Unit]
Description=Online XFS Metadata Check for All Filesystems
ConditionACPower=true
Documentation=man:xfs_scrub_all(8)

[Service]
Type=oneshot
Environment=SERVICE_MODE=1
ExecStart=@sbindir@/xfs_scrub_all
SyslogIdentifier=xfs_scrub_all
