# src/fs/gfs2/gfs2-utils/gfs2/scripts/82-gfs2-withdraw.rules
SUBSYSTEM=="gfs2", ACTION=="offline", RUN+="/bin/sh /usr/sbin/gfs2_withdraw_helper"

